# deploy/conduit-server.service
[Unit]
Description=conduit rendezvous + embedded TURN
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=conduit
Group=conduit
# /etc/conduit/server.env should define CONDUIT_TURN_PUBLIC_IP and
# CONDUIT_TURN_SECRET (and any other CONDUIT_SERVER_* overrides).
EnvironmentFile=/etc/conduit/server.env
ExecStart=/usr/local/bin/conduit-server \
  -addr 127.0.0.1:8080 \
  -trust-xff \
  -turn-embed \
  -turn-public-ip ${CONDUIT_TURN_PUBLIC_IP} \
  -turn-secret ${CONDUIT_TURN_SECRET}
# Bind 3478 without running as root.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
Restart=on-failure
RestartSec=5s

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
RestrictSUIDSGID=true
RestrictRealtime=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
