# deploy/conduit-server.service
[Unit]
Description=conduit rendezvous / signaling server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=conduit
Group=conduit
# /etc/conduit/server.env should define the CONDUIT_SERVER_* settings for the
# chosen TURN provider, e.g. CONDUIT_SERVER_CLOUDFLARE_TURN_KEY_ID and
# CONDUIT_SERVER_CLOUDFLARE_TURN_API_TOKEN for Cloudflare Realtime TURN.
EnvironmentFile=/etc/conduit/server.env
ExecStart=/usr/local/bin/conduit-server -addr 127.0.0.1:8080 -trust-xff
Restart=on-failure
RestartSec=5s

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
RestrictSUIDSGID=true
RestrictRealtime=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
